# CI/mssql-backup.timer
[Unit]
Description=Run SQL Server -> Azure Blob backup daily

[Timer]
OnCalendar=*-*-* 03:30:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
